# dockerfiles/app/dev/php-fpm/crontab
* * * * * /usr/local/bin/php /var/www/html/bin/console app:user_dish_choice_planned:set
